# rc.d/dynamic_motd.in
#!/bin/sh
# SPDX-License-Identifier: BSD-2-Clause
#
# Copyright (c) 2022 Steven Stallion <sstallion@gmail.com>
#
# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions
# are met:
# 1. Redistributions of source code must retain the above copyright
#    notice, this list of conditions and the following disclaimer.
# 2. Redistributions in binary form must reproduce the above copyright
#    notice, this list of conditions and the following disclaimer in the
#    documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE AUTHOR AND CONTRIBUTORS "AS IS" AND
# ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
# IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE
# ARE DISCLAIMED.  IN NO EVENT SHALL THE AUTHOR OR CONTRIBUTORS BE LIABLE
# FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL
# DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS
# OR SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS INTERRUPTION)
# HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT
# LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY
# OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF
# SUCH DAMAGE.

# PROVIDE: dynamic_motd
# REQUIRE: mountcritremove FILESYSTEMS
# BEFORE:  LOGIN

. /etc/rc.subr

name="dynamic_motd"
desc="Dynamic /var/run/motd"
rcvar="${name}"

set_rcvar dynamic_motd "NO"
set_rcvar motd_script "@@PREFIX@@/etc/rc.motd"

load_rc_config $name

command="@@PREFIX@@/libexec/${name}"
command_args="${motd_script}"
command_interpreter="/bin/sh"
start_precmd="${name}_prestart"
required_files="${motd_script}"

dynamic_motd_prestart()
{
	! checkyesno update_motd ||
	    err 1 "Set update_motd to NO in /etc/rc.conf."
}

run_rc_command "$1"
